# consul/server/supervisor.d/consul.ini
[program:consul]
command=/bin/consul agent -server -bootstrap -ui-dir=/ui -data-dir=/data -config-dir=/etc/consul.d -dns-port=53
autostart=true
autorestart=false
stdout_logfile=/tmp/%(program_name)s.log
stderr_logfile=/tmp/%(program_name)s.log
