# rails/supervisor.d/rails.ini
[program:rails]
command=/bin/sh -c "cd app; bundle install --system && rm -rf tmp/pids/server.pid && bin/rails s -b 0.0.0.0"
autostart=true
autorestart=false
stdout_logfile=/tmp/%(program_name)s.log
stderr_logfile=/tmp/%(program_name)s.log
